# packaging/linux/magicpad-gestures.service
[Unit]
Description=MagicPad Companion multi-finger gesture daemon
Documentation=https://github.com/imcmurray/MagicPad3
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
# Installed path may be /usr/local/bin or ~/.local/bin — the app rewrites this.
# `sg input` applies the input group without requiring a full re-login.
ExecStart=/usr/bin/sg input -c '/usr/local/bin/magicpad-companion --gestures'
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
